# services/init_start.service
[Unit]
Description=Initialize mDVR project setup
After=network.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/opt/mdvr
ExecStart=/opt/mdvr/venv/bin/python /opt/mdvr/init_project.py

[Install]
WantedBy=multi-user.target
