# services/mdvr_upload.timer
[Unit]
Description=Upload video to ftp

[Timer]
Unit=mdvr_upload.service
OnBootSec=3min
OnUnitActiveSec=3min

[Install]
WantedBy=timers.target
