# services/mdvr_rs.timer
[Unit]
Description=Run mdvr rs engine

[Timer]
Unit=mdvr_rs.service
OnBootSec=1min
OnUnitActiveSec=30min

[Install]
WantedBy=timers.target
