# services/mdvr_web.timer
[Unit]
Description=Timer for mDVR web server (on boot)

[Timer]
Unit=mdvr_web.service
OnBootSec=1min

[Install]
WantedBy=timers.target
