# services/mdvr_web.timer
[Unit]
Description=Run mdvr web

[Timer]
Unit=mdvr_web.service
OnBootSec=1min

[Install]
WantedBy=timers.target
