# services/mdvr.timer
[Unit]
Description=Run mdvr engine

[Timer]
Unit=mdvr.service
OnBootSec=1min

[Install]
WantedBy=timers.target
