# services/mdvr.timer
[Unit]
Description=Timer for mDVR recording service (on boot)

[Timer]
Unit=mdvr.service
OnBootSec=1min

[Install]
WantedBy=timers.target
